# anon_dumps/admin-master-deb4694b0e9110ffcf84a42f70874a6e152c0b32/need/supervisor/commands.conf
[program:commands]
directory=/var/www/html
command=/usr/bin/php index.php --task=commands
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
user=root
