# anon_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/blacklist_clear.conf
[program:blacklist_clear]
directory=/home/sendmail
command=/usr/bin/php artisan blacklist:clear
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0