# anon_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/clients_cache_delete.conf
[program:clients_cache_delete]
directory=/home/sendmail
command=/usr/bin/php artisan clients:cache:delete
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0