# anonfile_dumps/import-master-ac16d180c391fce7a644f6c2a30fc3cfb37451f6/need/supervisor/autoclear_network_archive.conf
[program:autoclear_network_archive]
directory=/root/import_1
command=/usr/bin/nice -n 19 /usr/bin/ionice -c2 -n7 /usr/bin/php index.php --task=autoclear:networkArchive
environment=KOHANA_ENV="production"
autostart=true
autorestart=true
numprocs=1
startsecs=0
user=root