# anonfile_dumps/cadmin-master-b2675af7f27c05513f1fd8374ee7bc35a058f18f/need/supervisor/update_cachenet.conf
[program:update_cachenet]
directory=/var/www/html
command=/usr/bin/nice -n 19 /usr/bin/ionice -c2 -n7 php index.php --task=update:cacheNet
environment=KOHANA_ENV="production"
process_name=%(program_name)s
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
