# anonfile_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/blacklist_queue_clients_online.conf
[program:blacklist_queue_clients_online]
directory=/home/sendmail
command=/usr/bin/php artisan blacklist:queue:clients:online
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0