# anonfile_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/blacklist_queue_check.conf
[program:blacklist_queue_check]
directory=/home/sendmail
command=/usr/bin/php artisan blacklist:queue:check
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0