# anonfile_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/emails_list_resolve_downloader.conf
[program:emails_list_resolve_downloader]
directory=/home/sendmail
command=/usr/bin/php artisan emails:list:resolve:downloader
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
user=root
startsecs=0