# anonfile_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/tasks_cache_email_list.conf
[program:tasks_cache_email_list]
directory=/home/sendmail
command=/usr/bin/php artisan tasks:cache:email:list
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0