# anonfile_dumps/sendmail-master-0a343a19f4f48dd8efd6c052c092fd5feec916ad/need/supervisor/tasks_cache_result.conf
[program:tasks_result_cache]
directory=/home/sendmail
command=/usr/bin/php artisan tasks:cache:result
numprocs=1
autostart=true
autorestart=true
user=root
startsecs=0