# anonfile_dumps/admin-master-deb4694b0e9110ffcf84a42f70874a6e152c0b32/need/supervisor/clients_filter.conf
[program:clients_filter]
directory=/var/www/html
command=/usr/bin/php index.php --task=clientsFilter
environment=KOHANA_ENV="production"
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
user=root
