# anonfile_dumps/admin-master-deb4694b0e9110ffcf84a42f70874a6e152c0b32/need/supervisor/nat.conf
[program:nat]
directory=/var/www/html
command=/usr/bin/php index.php --task=nat
environment=KOHANA_ENV="production"
autostart=true
autorestart=true
numprocs=1
user=root
stopsignal=KILL
startsecs=0