# anonfile_dumps/admin-master-deb4694b0e9110ffcf84a42f70874a6e152c0b32/need/supervisor/sphinx.conf
[program:admin_sphinx]
command=/usr/bin/searchd --config /etc/sphinxsearch/admin/sphinx.conf --nodetach --iostats --cpustats
autostart=true
autorestart=true
numprocs=1
user=root
